# conf/pactlers.service
[Unit]
Description=pactlers

[Service]
Type=simple
ExecStart=%h/.cargo/bin/pactlers
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
